# doc/robot_send_task.conf
[program:robot_send_task]
directory = /root/other_bid_project/Muto
command = /root/other_bid_project/ToolOfAnalysls/venv/bin/python -u /root/other_bid_project/Muto/plantform/robot_send_task.py
autostart = true
autorestart = true
# user = loanward

# 把stderr重定向到stdout，默认false;
redirect_stderr = true
# 标准日志输出;
stdout_logfile=/data/logs/robot_send_task.log
# 错误日志输出;
stderr_logfile=/data/logs/err-robot_send_task.log
# 标准日志文件大小，默认50MB;
stdout_logfile_maxbytes = 200MB
# 标准日志文件备份数;
# stdout_logfile_backups = 20
